# docker/supervisord.conf
[supervisord]
nodaemon=true
user=root
logfile=/dev/stdout
logfile_maxbytes=0
pidfile=/run/supervisord.pid

[program:php-fpm]
command=php-fpm -F
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
autorestart=true

[program:nginx]
command=nginx -g "daemon off;"
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
autorestart=true

; Dedicated queue worker — runs the heavy AI review jobs OUT of the web
; process. --memory=150 makes the worker exit (supervisor then restarts it)
; before it can approach the 512MB box limit; --max-time recycles it hourly
; to release slow leaks. --tries/--backoff mirror the job's own retry policy.
; Single worker = one review at a time = predictable memory.
[program:queue-worker]
command=php /var/www/html/artisan queue:work --queue=default --sleep=3 --tries=3 --backoff=60,180,600 --timeout=120 --memory=150 --max-time=3600
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
autorestart=true
stopwaitsecs=130
